# bin/service/dmarc_syslog.timer
[Unit]
Description=Extracts DMARC data from the email inbox

[Timer]
Unit=dmarc_syslog.service
OnCalendar=*:0/30

[Install]
WantedBy=timers.target
