# bin/service/dmarc_syslog.service
[Unit]
Description=Extracts DMARC data from the email inbox

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/dmarc_syslog
ExecStart=/usr/bin/python3 start.py

[Install]
WantedBy=multi-user.target
