# service/dmarc_syslog.timer
[Unit]
Description=Extracts DMARC data from the email inbox

[Timer]
Unit=dmarc_syslog.service
OnCalendar=*:0/5

[Install]
WantedBy=timers.target
